# e2e/conf/crontab
0 */2 * * * /usr/local/bin/python /src/main.py >> /src/log.log