# e2e/conf/crontab
* */2 * * * python /src/main.py >> /src/log.log